# bin/dashboard.service
[Unit]
Description=Dashboard
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/$user/Fiber-optic-project
ExecStart="/run.sh"
Restart=always

[Install]
WantedBy=multi-user.target
